# clamav/clamd.service
[Unit]
Description=Clamav
[Service]
After=freshclam
ExecStart=/usr/local/sbin/clamd
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
